# 01-servers/07-cluster-design/02-postgres-ha-cluster/ha-cluster-ansible/roles/etcd/templates/etcd.service.j2
[Unit]
Description=etcd
After=network.target

[Service]
Type=notify
User=etcd

ExecStart=/usr/local/bin/etcd \
  --name {{ inventory_hostname }} \
  --data-dir {{ etcd_data_dir }} \
  --initial-advertise-peer-urls http://{{ ansible_host }}:{{ etcd_peer_port }} \
  --listen-peer-urls http://{{ ansible_host }}:{{ etcd_peer_port }} \
  --listen-client-urls http://{{ ansible_host }}:{{ etcd_client_port }},http://127.0.0.1:{{ etcd_client_port }} \
  --advertise-client-urls http://{{ ansible_host }}:{{ etcd_client_port }} \
  --initial-cluster {{ etcd_initial_cluster }} \
  --initial-cluster-token {{ etcd_cluster_token }} \
  --initial-cluster-state {{ etcd_cluster_state | default('new') }}

Restart=always
RestartSec=5
LimitNOFILE=40000

[Install]
WantedBy=multi-user.target
